# koji-sweep-db/anacrontab
# Disable attempts to email job output anywhere
MAILTO=""

@daily 0 koji-sweep-db /usr/sbin/koji-sweep-db
